# projects/litcoin/litcoiin-miner.service
[Unit]
Description=Manteclaw Litcoiin Standalone Miner
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/projects/litcoin
Environment="LITCOIN_SEED_FILE=/root/.openclaw/workspace/.keys/wallet.seed"
Environment="FIREWORKS_API_KEY_FILE=/root/.openclaw/workspace/.keys/.env"
Environment="PYTHONUNBUFFERED=1"
Environment="NVIDIA_API_KEY_FILE=/root/.openclaw/workspace/.keys/.env"
ExecStart=/root/.openclaw/workspace/projects/litcoin/venv/bin/python3 /root/.openclaw/workspace/projects/litcoin/standalone-miner.py 0 3
Restart=always
RestartSec=10
StandardOutput=append:/root/.openclaw/workspace/projects/litcoin/miner_service.log
StandardError=append:/root/.openclaw/workspace/projects/litcoin/miner_service.log

[Install]
WantedBy=default.target
